# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:auth-proxy]
process_name=%(program_name)s
command=node /app/auth-proxy/index.cjs
environment=AUTH_PROXY_PORT="%(ENV_PORT)s",OPENCLAW_PORT="8082",FILEBROWSER_PORT="8081",PAGES_API_PORT="8083",PROXY_USER="%(ENV_PROXY_USER)s",PROXY_PASSWORD="%(ENV_PROXY_PASSWORD)s",AUTH_SECRET="%(ENV_AUTH_SECRET)s"
autostart=true
autorestart=true
user=root
numprocs=1
startsecs=3
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
priority=100

[program:pages-api]
process_name=%(program_name)s
command=node /app/backend/server.cjs
environment=PAGES_API_PORT="8083",OPENCLAW_WORKSPACE_DIR="%(ENV_OPENCLAW_WORKSPACE_DIR)s"
autostart=true
autorestart=true
user=root
numprocs=1
startsecs=2
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
priority=110

[program:openclaw]
process_name=%(program_name)s
command=node /app/dist/index.js gateway --allow-unconfigured --port 8082 --bind loopback
autostart=true
autorestart=true
user=root
numprocs=1
startsecs=10
startretries=5
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stopsignal=TERM
priority=200

[program:filebrowser]
process_name=%(program_name)s
command=filebrowser --baseurl /files --address 127.0.0.1 --port 8081 --root /data --database /data/.openclaw/filebrowser.db --noauth
autostart=true
autorestart=true
user=root
numprocs=1
startsecs=3
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
priority=150
